# launchagents/com.jopcmelo.force-audio-input.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jopcmelo.force-audio-input</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__DOTFILES_DIR__/scripts/force-audio-input.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/force-audio-input.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/force-audio-input.err</string>
</dict>
</plist>
